# gen8/gen8-update-split-cidrs.timer
[Unit]
Description=Weekly refresh of gen8 split VPN CIDR table

[Timer]
OnCalendar=Sun 04:20
RandomizedDelaySec=30min
Persistent=true

[Install]
WantedBy=timers.target
